# supervisord.conf
[supervisord]
nodaemon=true

[program:crontabui]
command=crontab-ui
stderr_logfile = /var/log/crontabui-stderr.log
stdout_logfile = /var/log/crontabui-stdout.log